# measurements/demo-workload
#!/sbin/openrc-run

command="/srv/workload.sh"
command_background="yes"
pidfile="/tmp/workload.pid"

